# com.daily-audio-pipeline.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daily-audio-pipeline.dashboard</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/anaconda3/bin/python3</string>
        <string>/Users/dongyue/Documents/daily-audio-pipeline/dashboard.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/dongyue/Documents/daily-audio-pipeline</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/dongyue/Documents/daily-audio-pipeline/output/dashboard.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/dongyue/Documents/daily-audio-pipeline/output/dashboard.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>TTS_ENGINE</key>
        <string>edge-tts</string>
    </dict>
</dict>
</plist>
